# setup/bluetooth_midi_setup/system/lib/systemd/system/midi.service
[Unit]
Description=Initial USB MIDI connect

[Service]
ExecStart=/usr/local/bin/remidi_connect.rb

[Install]
WantedBy=multi-user.target
